# 2017/solutions/poetwanna.be/config/supervisor/client-ec2/conf.d/qa_nips.conf
[program:qa_nips]
command = /webapps/qa_nips/qa_nips/django/bin/gunicorn_start           ; Command to start app
user = qanips                                                          ; User to run as
stdout_logfile = /webapps/qa_nips/logs/gunicorn_supervisor.log         ; Where to write log messages
redirect_stderr = true                                                 ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                        ; Set UTF-8 as default encoding
